# Ansible/files/rc.local
#!/bin/sh -e

ifconfig wlan0 promisc
exit 0
